# projekteserver/systemd/examibur-backup.service
[Unit]
Description=Backup all docker volumes

[Service]
Type=simple
ExecStart=/opt/examibur/projekteserver/backup/backup.sh

